# deploy/systemd/caddy.service.d/boke-video.conf.example
[Service]
EnvironmentFile=/etc/boke-video/caddy.env
